# scripts/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

#java -cp /home/ec2-user/CNV_PARA_O_20/ WebServer
(cd /home/ec2-user/CNV_PARA_O_20/ && make run_httpserver)
(wget -q -O - http://169.254.169.254/latest/meta-data/instance-id) >> /home/ec2-user/id.txt
